# base/cron.d/1logrotate
02 00 * * * root logrotate -v /etc/logrotate.conf