# base/cron.d/2tilecache
03 00 * * * root /opt/odtiles/tilecache_rotate.sh